# MESH/safranyHope/soil/systemd/soil-sample.timer
[Unit]
Description=Run soil sampler hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
